# deploy/ursus.service
# Systemd unit для автозапуска Ursus Insight SIEM
#
# Установка:
#   sudo cp ursus.service /etc/systemd/system/
#   sudo systemctl daemon-reload
#   sudo systemctl enable ursus
#   sudo systemctl start ursus
#
# Просмотр логов:
#   sudo journalctl -u ursus -f

[Unit]
Description=Ursus Insight SIEM
After=network.target

[Service]
Type=simple
User=www-data
WorkingDirectory=/opt/ursus
ExecStart=/usr/bin/python3 /opt/ursus/Main.py --no-demo
Restart=always
RestartSec=5

# Обязательно смените пароли!
Environment=URSUS_USER=admin
Environment=URSUS_PASS=your-strong-password-here
Environment=URSUS_AGENT_KEY=your-agent-key-here
Environment=URSUS_SECRET=your-secret-key-here

# Hardening
NoNewPrivileges=yes
PrivateTmp=yes

[Install]
WantedBy=multi-user.target
